# mysql/supervisord.conf
[supervisord]
nodeamon=true

[program:mysqld]
command=/usr/bin/mysqld_safe
